# systemd/anolis-provider-bread.service
[Unit]
Description=Anolis Provider — Bread (CRUMBS Bus)
After=network.target
PartOf=anolis-runtime.service

[Service]
Type=simple
User=anolis
Group=anolis
SupplementaryGroups=i2c gpio dialout
ExecStart=/opt/anolis/bin/anolis-provider-bread --config /opt/anolis/config/providers/bread.yaml
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
